# scripts/bigbox.service
[Unit]
Description=bigbox pentesting firmware
After=multi-user.target network.target
ConditionPathExists=/opt/bigbox/bigbox/__main__.py

[Service]
Type=simple
# Root so GPIO + raw sockets + reboot/poweroff actions just work. Tighten
# later with capabilities (CAP_NET_RAW, CAP_NET_ADMIN, CAP_SYS_TIME) if you
# want to drop to a non-root user.
User=root
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0
WorkingDirectory=/opt/bigbox

# Launch a single-purpose X server on vt1 and run bigbox as its only client.
ExecStartPre=-/usr/bin/amixer -c 1 sset "PCM" 100% unmute
ExecStart=/usr/bin/xinit /opt/bigbox/.venv/bin/python -m bigbox -- :0 vt1 -nocursor -nolisten tcp

TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
StandardInput=tty
StandardOutput=journal
StandardError=journal

Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
